# deploy/uptrend-sideways-paper-trader/uptrend-sideways-paper-trader.service
[Unit]
Description=Uptrend Sideways Demo Paper Trader
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/hacker/apps/swing-trading-app
ExecStart=/usr/bin/docker compose -f deploy/matsya-setup/docker-compose.yml --profile manual run --rm uptrend-sideways-paper-trader

[Install]
WantedBy=multi-user.target
