# deploy/matsya-candle-mirror/matsya-candle-mirror.cron
# Matsya primary-to-swing candle mirror (single writer: new VPS)
*/2 * * * * /home/hacker/.local/bin/matsya-candle-mirror-run
