# deploy/matsya-candle-mirror/matsya-candle-mirror.cron
# Matsya primary-to-swing candle mirror (single writer: new VPS)
30 6 * * * /home/hacker/.local/bin/matsya-candle-mirror-run
