# deploy/matsya-candle-mirror/matsya-candle-mirror.timer
[Unit]
Description=Continuously mirror Matsya Dhan daily candles

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=true
Unit=matsya-candle-mirror.service

[Install]
WantedBy=timers.target
