# deploy/matsya-candle-mirror/matsya-candle-mirror.service
[Unit]
Description=Mirror Matsya Dhan daily candles from the primary VPS
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=hacker
Group=hacker
EnvironmentFile=/etc/matsya-candle-mirror.env
ExecStart=/usr/local/sbin/matsya-candle-mirror-sync
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
